# templates/borgbackup.timer.j2
[Unit]
Description=Borgbackup timer

[Timer]
OnCalendar={{ borgbackup_config.calendar_spec|default("*-*-* 2:00:00") }}
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
