# templates/borgbackup-tasks.service.j2
[Unit]
Description=Pre backup tasks for 'borgbackup'

[Service]
User={{ borgbackup_user }}
ExecStart=/bin/true
Type=oneshot
