# templates/borgbackup-tasks.service.j2
[Unit]
Description=Backup tasks collector

[Service]
ExecStart=/bin/sh -c "! /bin/systemctl is-failed 'borgbackup-tasks@*'"
Type=oneshot
