# templates/borgbackup-tasks.service.j2
[Unit]
Description=Backup {{ service }}
Before=borgbackup.service

[Service]
User={{ backup_user }}
ExecStart={{ backup_command }}
Type=oneshot

[Install]
WantedBy=borgbackup.service
