# templates/borgbackup.timer.j2
[Unit]
Description=Borgbackup timer

[Timer]
OnCalendar={{ borgbackup_calendar_spec }}
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
